# systemd/system/harbour-storeman-installer.timer
[Unit]
Description=timer unit to trigger the unit for installing harbour-storeman

[Timer]
# Starts the corresponding service unit two to three seconds after activation:
OnActiveSec=2
AccuracySec=1
# Allows for retriggering this timer unit without rebooting, but makes debugging harder:
RemainAfterElapse=false
Unit=harbour-storeman-installer.service  # Should not, but seems to be necessary!

