# systemd/system/harbour-storeman-installer.timer
[Unit]
Description=Trigger service unit for installing Storeman after a second

[Timer]
OnActiveSec=1
