# systemd/system/harbour-storeman-installer.service
[Unit]
Description=oneshot unit for installing Storeman

[Service]
Type=oneshot
ExecStartPre=/usr/bin/pkcon -pn repo-set-data harbour-storeman-obs refresh-now true
ExecStart=/usr/bin/pkcon -pny install harbour-storeman
